# services/thyra.conf
; /etc/supervisor/conf.d/thyra.conf
; Thyra — gestion des processus par Supervisord
; Pas de Docker, pas de systemd complexe : supervisord gère tout.

[unix_http_server]
file=/run/supervisor.sock
chmod=0770
chown=root:www-data

[supervisord]
logfile=/var/log/thyra/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/run/supervisord.pid
nodaemon=false
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

; ── Gunicorn / Flask server ───────────────────────────────────────────────
[program:thyra-server]
command=/opt/thyra/venv/bin/gunicorn
        --workers 2
        --worker-class sync
        --bind 127.0.0.1:5000
        --timeout 120
        --keep-alive 5
        --access-logfile /var/log/thyra/gunicorn_access.log
        --error-logfile  /var/log/thyra/gunicorn_error.log
        --log-level info
        server:app
directory=/opt/thyra/app
user=thyra
environment=THYRA_HOME="/opt/thyra",THYRA_USER="thyra"
autostart=true
autorestart=true
startsecs=3
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/thyra/server_stdout.log
stderr_logfile=/var/log/thyra/server_stderr.log

; ── Viewer ───────────────────────────────────────────────────────────────
[program:thyra-viewer]
command=/opt/thyra/venv/bin/python3 /opt/thyra/app/viewer.py
directory=/opt/thyra/app
user=thyra
environment=THYRA_HOME="/opt/thyra",THYRA_USER="thyra",DISPLAY=":0",XAUTHORITY="/home/thyra/.Xauthority"
autostart=true
autorestart=true
startsecs=8
startretries=10
; Démarre après le serveur
; (supervisor ne garantit pas l'ordre strict, le viewer boucle sur l'API)
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/thyra/viewer_stdout.log
stderr_logfile=/var/log/thyra/viewer_stderr.log

; ── X11 / Openbox minimal (pour Chromium kiosk pages web) ─────────────────
[program:thyra-x11]
command=/bin/bash -c "rm -f /tmp/.X0-lock; /usr/bin/Xorg :0 -nocursor vt1"
user=root
autostart=true
autorestart=true
startsecs=3
priority=10
stdout_logfile=/var/log/thyra/x11_stdout.log
stderr_logfile=/var/log/thyra/x11_stderr.log

; ── Watchdog : vérifie que l'écran ne s'éteint pas ────────────────────────
[program:thyra-watchdog]
command=/opt/thyra/scripts/watchdog.sh
user=root
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/thyra/watchdog.log
stderr_logfile=/var/log/thyra/watchdog.log

; ── Portail captif (désactivé par défaut, activé via settings) ───────────
[program:thyra-ap]
command=/opt/thyra/scripts/ap_manager.sh run
user=root
autostart=false
autorestart=false
stdout_logfile=/var/log/thyra/ap.log
stderr_logfile=/var/log/thyra/ap.log
